# config/var-lib-docker.mount
[Unit]
Description=mount persistent drive to /var/lib/docker
[Mount]
What=/dev/disk/by-id/google-stable-docker
Where=/var/lib/docker
Type=ext4
