# docker/fedora-ansible/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:sshd]
command=/usr/sbin/sshd -eD
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/sshd.log
